# bot/crontab
*/2 * * * * API_ID='your_api_id' API_HASH='your_api_hash' PHONE='your_phone' DB_URL='postgresql://user:qwe123@postgre:5432/messagesdb' /usr/local/bin/python3 /bot/main.py >> /bot/cron.log 2>&1
